# etc/s6-overlay/s6-rc.d/ssr/run
#!/usr/bin/env bash
php "$APP_BASE_DIR/artisan" inertia:start-ssr --runtime=bun
